# etc/supervisor/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/var/log/supervisor/supervisord.log 
pidfile=/var/run/supervisord.pid 
childlogdir=/var/log/supervisor            

[unix_http_server]
file=/var/run/supervisor.sock   
chmod=0700                       

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock 

[program:node-red-admin]
command= node-red -s /rules/bots/admin/.node-red/settings.js -u /rules/bots/admin

[program:node-red-main]
command= node-red -s /rules/bots/main/.node-red/settings.js -u /rules/bots/main

[program:node-red-bot1]
command= node-red -s /rules/bots/bot1/.node-red/settings.js -u /rules/bots/bot1

[program:node-red-bot2]
command= node-red -s /rules/bots/bot2/.node-red/settings.js -u /rules/bots/bot2

[program:node-red-bot3]
command= node-red -s /rules/bots/bot3/.node-red/settings.js -u /rules/bots/bot3

[program:node-red-bot4]
command= node-red -s /rules/bots/bot4/.node-red/settings.js -u /rules/bots/bot4

[program:node-red-bot5]
command= node-red -s /rules/bots/bot5/.node-red/settings.js -u /rules/bots/bot5

[program:node-red-bot6]
command= node-red -s /rules/bots/bot6/.node-red/settings.js -u /rules/bots/bot6

[program:node-red-bot7]
command= node-red -s /rules/bots/bot7/.node-red/settings.js -u /rules/bots/bot7

[program:node-red-bot8]
command= node-red -s /rules/bots/bot8/.node-red/settings.js -u /rules/bots/bot8

[program:node-red-bot9]
command= node-red -s /rules/bots/bot9/.node-red/settings.js -u /rules/bots/bot9

[program:node-red-bot10]
command= node-red -s /rules/bots/bot10/.node-red/settings.js -u /rules/bots/bot10
